# files/e3372-watchdog.timer
[Unit]
Description=e3372 watchdog каждые 30с
[Timer]
OnBootSec=20
OnUnitActiveSec=30
AccuracySec=5s
[Install]
WantedBy=timers.target
