# files/e3372-watchdog.service
[Unit]
Description=e3372 watchdog (routing + modem health + DNS safety)
After=systemd-networkd.service
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/e3372-watchdog.sh
